# .docker/base/etc/service/haproxy/run
#!/bin/sh

exec /usr/sbin/haproxy -f /etc/haproxy/haproxy.conf -sf
